# etc/powerman.service.in
[Unit]
Description=PowerMan
After=syslog.target network.target

[Service]
Environment=SHELL=/bin/sh
PrivateTmp=yes
User=@RUN_AS_USER@
Group=@RUN_AS_GROUP@
ExecStart=@X_SBINDIR@/powermand -f

[Install]
WantedBy=multi-user.target
