# scripts/powerman.service.in
[Unit]
Description=PowerMan
After=syslog.target network.target

[Service]
Type=forking
PrivateTmp=yes
User=@RUN_AS_USER@
Group=@RUN_AS_GROUP@
ExecStart=/usr/sbin/powermand
PIDFile=/var/run/powerman/powermand.pid

[Install]
WantedBy=multi-user.target
